# scripts/crontab
# m h  dom mon dow   cmd

# 30 14 * * * /home/appvilla/graphhopper/start.sh >> /home/appvilla/graphhopper/start-crontab.log 2>&1

* * * * * env > /tmp/env.output






